# rootfs/etc/services.d/60-ibc-logs/run
#!/bin/sh

exec tail -f /root/ibc/logs/ibc*
