# rootfs/etc/services.d/90-socat/run
#!/usr/bin/env bash

exec socat TCP-LISTEN:4002,fork TCP:127.0.0.1:4001
