# rootfs/etc/services.d/10-xvfb/run
#!/bin/sh

export DISPLAY=":0"

exec /usr/bin/Xvfb "$DISPLAY" -ac -screen 0 1024x768x24
