# rootfs/etc/services.d/60-ibcontroller-logs/run
#!/bin/sh

exec tail -f /root/IBController/Logs/ibc*
